# supervisord.conf
[supervisord]
nodaemon=true

[program:web]
command=python manage.py runserver 0.0.0.0:8000
autostart=true
autorestart=true

[program:worker]
command=python blog/consumer.py
autostart=true
autorestart=true
